# marwahicloudkichan_update-main/hdos.service
[Unit]
Description=HDOS - Hotel Digital Operating System
After=network.target postgresql.service
Requires=postgresql.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/hdos
Environment=NODE_ENV=production
EnvironmentFile=/var/www/hdos/.env
ExecStart=/usr/bin/node /var/www/hdos/dist/index.cjs
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hdos

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ReadWritePaths=/var/www/hdos/uploads
ProtectHome=yes
PrivateDevices=yes

[Install]
WantedBy=multi-user.target